# contrib/init/reecored.init
#!/bin/bash
#
#  relianced The Reliance core server.
#
#
# chkconfig: 345 80 20
# description: relianced
# processname: relianced
#

# Source function library.
. /etc/init.d/functions

# you can override defaults in /etc/sysconfig/relianced, see below
if [ -f /etc/sysconfig/relianced ]; then
        . /etc/sysconfig/relianced
fi

RETVAL=0

prog=relianced
# you can override the lockfile via BITCOIND_LOCKFILE in /etc/sysconfig/relianced
lockfile=${BITCOIND_LOCKFILE-/var/lock/subsys/relianced}

# relianced defaults to /usr/bin/relianced, override with BITCOIND_BIN
bitcoind=${BITCOIND_BIN-/usr/bin/relianced}

# relianced opts default to -disablewallet, override with BITCOIND_OPTS
bitcoind_opts=${BITCOIND_OPTS}

start() {
    echo -n $"Starting $prog: "
    daemon $DAEMONOPTS $bitcoind $bitcoind_opts
    RETVAL=$?
    echo
    [ $RETVAL -eq 0 ] && touch $lockfile
    return $RETVAL
}

stop() {
    echo -n $"Stopping $prog: "
    killproc $prog
    RETVAL=$?
    echo
    [ $RETVAL -eq 0 ] && rm -f $lockfile
    return $RETVAL
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    status)
        status $prog
        ;;
    restart)
        stop
        start
        ;;
    *)
        echo "Usage: service $prog {start|stop|status|restart}"
        exit 1
        ;;
esac
